# bash/etc/user.cron
# m h dom mon dow	command
5 * * * *	exec ~/httpsdocs/bg-process/run.sh check &
