# etc/Server/aconfmgr/calamares/etc/local.d/branding.start
#!/bin/sh

neofetch --no_config >| /etc/issue

